# conf/supervisor.conf
# -*- conf -*-

[program:uss]
command=/srv/www/uss/conf/startup
user=ubuntu
autostart=true
autorestart=true
stderr_logfile=/srv/www/logs/uss/supervisor_err.log
