# conf/supervisor.conf
# -*- conf -*-

[program:uss]
command=/srv/sites/uss/conf/startup
user=ubuntu
autostart=true
autorestart=true
stderr_logfile=/srv/logs/uss/supervisor_err.log
